# deployment/fastapi.socket
[Unit]
Description=fastapi socket

[Socket]
ListenStream=/run/fastapi.sock

[Install]
WantedBy=sockets.target
